# build/root_overlay/var/service/rpcbind/run
#!/bin/bash

exec 2>&1

# -s
exec /sbin/rpcbind -s -d


