# build/root_overlay/var/service/nginx/run
#!/bin/sh

#mkdir -p /var/tmp/nginx
#chown nginx:nginx /var/tmp/nginx

exec /usr/sbin/nginx

sleep 10

