# build/root_overlay/var/service/smbd/run
#!/bin/sh

mkdir -p /var/run/samba/

exec /usr/sbin/smbd -F

