# build/root_overlay/var/service/nmbd/run
#!/bin/sh

mkdir -p /var/run/samba/

exec /usr/sbin/nmbd -F -S

